# supervisord.conf
[supervisord]
nodaemon=true

[program:script]
command=/usr/local/bin/python3 revproxy.py
autorestart=true

[program:redis]
command=/usr/bin/redis-server
autorestart=true
